# templates/systemd.service.j2
{{ ansible_managed | comment }}
[Unit]
Description={{ tomcat_name }}
After=network.target

[Service]
Type=forking
PIDFile={{ tomcat_base }}/{{ tomcat_pid_file }}
User={{ tomcat_user }}
Group={{tomcat_group }}
Restart=on-failure
{% if tomcat_use_native_daemon
      and tomcat_version is version("7", ">=") %}
Environment=JAVA_HOME={{ tomcat_java_home }}
Environment=CATALINA_HOME={{ tomcat_home }}
Environment=CATALINA_BASE={{ tomcat_base }}
ExecStart={{ tomcat_base }}/bin/daemon.sh start
ExecStop={{ tomcat_base }}/bin/daemon.sh stop
{% else %}
ExecStart={{ tomcat_base }}/bin/startup.sh
ExecStop={{ tomcat_base }}/bin/shutdown.sh
{% endif %}
TimeoutStartSec=30
TimeoutStopSec=30
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
